# etc/supervisor/conf.d/dvwa.conf
[program:dvwa]
command=/usr/sbin/dvwa
user=root

